# deploy/supervisor/pocketapp-reverb.conf
; Keeps the app's own WebSocket broadcast server (Laravel Reverb) running
; permanently. Every user's chart subscribes to this (via Echo, see
; resources/js/echo.js) for live price ticks and other broadcasts — without
; it running, `broadcast()` calls succeed server-side but nothing is actually
; delivered to any connected browser.
;
; Install:
;   sudo cp deploy/supervisor/pocketapp-reverb.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start pocketapp-reverb:*
;
; Check status / tail logs:
;   sudo supervisorctl status pocketapp-reverb:*
;   tail -f storage/logs/reverb.log

[program:pocketapp-reverb]
process_name=%(program_name)s
command=php /var/www/pocketapp/artisan reverb:start
directory=/var/www/pocketapp
autostart=true
autorestart=true
startsecs=5
startretries=100
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/pocketapp/storage/logs/reverb.log
stopwaitsecs=10
