# deploy/supervisor/pocketapp-prices-sync.conf
; Keeps the live price/websocket listener running permanently.
;
; `php artisan prices:sync` is an infinite reconnect-loop command (see
; app/Console/Commands/SyncAssetPrices.php) — it is NOT meant to be run
; via cron/scheduler, and it does not survive a terminal/SSH session
; closing on its own. Without a supervisor keeping exactly one instance
; of this alive at all times, every asset will start reporting "offline"
; (and trading will be blocked) the moment the process dies for any
; reason — server reboot, deploy, OOM kill, an uncaught exception, etc.
;
; Install:
;   sudo cp deploy/supervisor/pocketapp-prices-sync.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start pocketapp-prices-sync:*
;
; Check status / tail logs:
;   sudo supervisorctl status pocketapp-prices-sync:*
;   tail -f storage/logs/prices-sync.log

[program:pocketapp-prices-sync]
process_name=%(program_name)s
command=php /var/www/pocketapp/artisan prices:sync
directory=/var/www/pocketapp
autostart=true
autorestart=true
startsecs=5
startretries=100
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/pocketapp/storage/logs/prices-sync.log
stopwaitsecs=10
