# deploy/supervisor/pocketapp-price-collector.conf
; Keeps the backend price collector (collector/index.js) running permanently.
; This is the ONE real connection to iqcent's price feed — a headless
; browser (Playwright/Chromium), not a plain WS/HTTP client, because iqcent
; sits behind Cloudflare and 403s any non-browser connection attempt (this
; was verified against both a PHP client with spoofed browser headers and a
; plain Node `ws` client — both rejected before the handshake completed).
;
; It relays ticks into PriceFeedService via PriceCollectorController, which
; broadcasts them out over Reverb (see pocketapp-reverb.conf — must also be
; running) to every connected user's chart. Without this process:
;   - no asset will ever report "online", so trading stays blocked
;   - no chart will move for anyone
;
; This replaced the old `php artisan prices:sync` daemon (deleted along with
; app/Console/Commands/SyncAssetPrices.php / app/Services/WebSocketListener.php)
; once it was confirmed that command's server-side WS connection gets the
; exact same Cloudflare 403 — it could never have worked.
;
; Install:
;   cd /var/www/pocketapp/collector && npm install   # also runs `playwright install`
;   sudo cp deploy/supervisor/pocketapp-price-collector.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start pocketapp-price-collector:*
;
; Required environment (set PRICE_COLLECTOR_SECRET to the same value as the
; Laravel app's .env, and point APP_URL at wherever this app is actually
; reachable from this machine):
;   APP_URL=http://127.0.0.1
;   PRICE_COLLECTOR_SECRET=<same value as pocketapp's .env PRICE_COLLECTOR_SECRET>
;
; Check status / tail logs:
;   sudo supervisorctl status pocketapp-price-collector:*
;   tail -f /var/www/pocketapp/storage/logs/price-collector.log

[program:pocketapp-price-collector]
process_name=%(program_name)s
command=node /var/www/pocketapp/collector/index.js
directory=/var/www/pocketapp/collector
; PRICE_COLLECTOR_SECRET is REQUIRED — collector/index.js refuses to start
; without it (see its own startup check). The docs above always said to set
; it, but this directive itself never actually carried it, so installing
; this file as-is silently produced a collector that never runs at all.
; deploy/setup-streaming.sh fills in __PRICE_COLLECTOR_SECRET__ from the
; Laravel app's .env automatically; if installing by hand instead, replace
; it yourself with that same value.
environment=APP_URL="http://127.0.0.1",PRICE_COLLECTOR_SECRET="__PRICE_COLLECTOR_SECRET__"
autostart=true
autorestart=true
startsecs=5
startretries=100
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/pocketapp/storage/logs/price-collector.log
stopwaitsecs=15
stopsignal=TERM
