# deploy/supervisor/pocketapp-queue-worker.conf
; Keeps trade settlement, cashback, notifications, and every other
; queued job (QUEUE_CONNECTION=database) actually processing.
;
; Trades are placed via EvaluateTrade::dispatch($trade)->delay(...) and
; ExpressTradeJob works the same way (see app/Http/Controllers/TradeController.php
; and app/Jobs/ExpressTradeJob.php) — win/loss settlement, payouts, and
; loss-cashback crediting all happen INSIDE these queued jobs. Without a
; worker permanently consuming the `jobs` table, dispatched jobs just sit
; there forever and trades never actually settle, no matter how correct
; the settlement logic itself is. (Verified in this sandbox: 41 jobs had
; queued up with zero worker ever running to consume them.)
;
; Install:
;   sudo cp deploy/supervisor/pocketapp-queue-worker.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start pocketapp-queue-worker:*
;
; After every deploy, reload workers so they pick up new code (queue:work
; boots once and keeps running the OLD code in memory otherwise):
;   php artisan queue:restart
;
; Check status / tail logs:
;   sudo supervisorctl status pocketapp-queue-worker:*
;   tail -f storage/logs/queue-worker.log

[program:pocketapp-queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/pocketapp/artisan queue:work database --sleep=1 --tries=3 --max-time=3600
directory=/var/www/pocketapp
autostart=true
autorestart=true
startsecs=5
startretries=100
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/pocketapp/storage/logs/queue-worker.log
stopwaitsecs=3600
