# deploy/supervisor/pocketapp-ticker-collector.conf
; Keeps the ticker collector (php artisan ticks:collect) running permanently,
; as a pool of processes — one per batch of ~10 assets, each driving its own
; headless Chrome instance via Panther.
;
; Why a pool instead of one process for the whole catalog: iqcent appears to
; cap how many SUBSCRIBE.TICK subscriptions a single WS connection actually
; honors — one connection subscribed to the full catalog only ever ticked for
; a small fraction of it, regardless of send pacing. Splitting the catalog
; into small batches, one Chrome client per batch, keeps each connection's
; subscription count low enough to stay reliable. See
; app/Http/Controllers/TickerController.php for the collection logic itself.
;
; Each process relays ticks straight into PriceFeedService in-process (no
; HTTP hop — unlike the old Node/Playwright collector this replaced), then
; broadcasts them out over Reverb (see pocketapp-reverb.conf — must also be
; running) to every connected user's chart. Without this process pool:
;   - no asset will ever report "online", so trading stays blocked
;   - no chart will move for anyone
;
; numprocs below MUST be >= ceil(total asset count / --size). At the time of
; writing this is ceil(158 / 10) = 16. If the asset catalog grows past
; numprocs * 10, the extra assets simply never get subscribed — bump
; numprocs (and re-run the install steps below) if Assets::count() grows.
; A process with no assets left to claim idles and checks back every 5
; minutes rather than erroring (see CollectTicks::handle) — safe to always
; run a couple more procs than currently needed as headroom.
;
; Install (deploy/setup.sh does all of this automatically — see README.md):
;   sudo cp deploy/supervisor/pocketapp-ticker-collector.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start pocketapp-ticker-collector:*
;
; Check status / tail logs:
;   sudo supervisorctl status pocketapp-ticker-collector:*
;   tail -f /var/www/pocketapp/storage/logs/ticker-collector.log
;
; Requires a Chrome/Chromium browser installed on this machine (Panther
; drives it, doesn't ship it) — see README.md's prerequisites. The matching
; chromedriver binary is already committed at drivers/chromedriver (Panther
; auto-discovers it there, since Supervisor's `directory=` below is the app
; root); Panther throws a clear "chromedriver binary not found" error if it
; ever can't find it. If the server's Chrome version drifts far enough from
; that binary's version to stop working, regenerate it:
;   cd /var/www/pocketapp && composer install && vendor/bin/bdi detect drivers

[program:pocketapp-ticker-collector]
process_name=%(program_name)s_%(process_num)s
command=php /var/www/pocketapp/artisan ticks:collect --batch=%(process_num)s --size=10
directory=/var/www/pocketapp
autostart=true
autorestart=true
startsecs=10
startretries=100
user=www-data
numprocs=16
numprocs_start=0
redirect_stderr=true
stdout_logfile=/var/www/pocketapp/storage/logs/ticker-collector.log
stopwaitsecs=15
stopsignal=TERM
